# puppet/shared_puppet/modules/fxmp_stunnel/templates/stunnel.erb
#!/bin/bash
#   /etc/rc.d/init.d/stunnel
#
# Starts one stunnel process for each config file
#
# chkconfig: - 85 15
# description: Stunnels
# Source function library.
if [ -f /etc/init.d/functions ]; then . /etc/init.d/functions; fi
if [ -f /lib/lsb/init-functions ]; then . /lib/lsb/init-functions; fi

CONFIGDIR="<%= conf_dir %>"

start() {
    echo -n $"Starting Stunnel $1: "
    pid=`ps -ef | grep stunnel | grep $1 | grep -v grep | awk {'print $2'}`
    if [ "$pid" != "" ]
    then
        echo " Already started"
        exit
    fi
    stunnel $1
    sleep 1
    pid=`ps -ef | grep stunnel | grep $1 | grep -v grep | awk {'print $2'}`
    [ "$pid" != "" ] && echo " OK" || echo " Failed"
}
stop() {
    echo -n $"Shutting down Stunnel $1: "
    pid=`ps -ef | grep stunnel | grep $1 | grep -v grep | awk {'print $2'}`
    if [ "$pid" != "" ]
    then
        kill -9 $pid > /dev/null
        sleep 1
        pid=`ps -ef | grep stunnel | grep $1 | grep -v grep | awk {'print $2'}`
        [ "$pid" == "" ] && echo " OK" || echo " Failed"
    else
        echo " Not running."
    fi
}

status(){
    pid=`ps -ef | grep stunnel | grep $1 | grep -v grep | awk {'print $2'}`
    if [ "$pid" != "" ]
    then
        echo "Stunnel $1 Running"
    else
        echo "Stunnel $1 Stopped"
    fi
}
########################################
##   MAIN
########################################

RETVAL=0

umask 002


case "$1" in
 start)
       if [ "$2" == "" ]
       then
            for i in $(ls $CONFIGDIR/stunnel.*.conf);
            do
                start $i
            done
        else
            if [ -f $CONFIGDIR/stunnel.$2.conf ]
            then
                start $CONFIGDIR/stunnel.$2.conf
            else
                echo "Config file $CONFIGDIR/stunnel.$2.conf not found."
            fi
        fi
       ;;
 stop)
    if [ "$2" == "" ]
       then
            for i in $(ls $CONFIGDIR/stunnel.*.conf);
            do
                stop $i
            done
        else
            stop $CONFIGDIR/stunnel.$2.conf
        fi
       ;;
 restart|reload)
       if [ "$2" == "" ]
       then
            for i in $(ls $CONFIGDIR/stunnel.*.conf);
            do
                stop $i
                start $i
            done
        else
            if [ -f $CONFIGDIR/stunnel.$2.conf ]
            then
                stop $CONFIGDIR/stunnel.$2.conf
                start $CONFIGDIR/stunnel.$2.conf
            else
                echo "Config file $CONFIGDIR/stunnel.$2.conf not found."
            fi
        fi
       ;;
 status)
       if [ "$2" == "" ]
       then
            for i in $(ls $CONFIGDIR/stunnel.*.conf);
            do
                status $i
            done
        else
            status $CONFIGDIR/stunnel.$2.conf
        fi
      ;;
 *)
       echo $"Usage: $0 {start|stop|restart|status}"
       exit 1
esac

exit $?
